# benchmarks/linux-node/faas.openrc
#!/sbin/openrc-run

command="node"
command_args="/bin/main.js $(sed "s/.*callback=\([^ ]*\).*/\1/" /proc/cmdline)"
pidfile="/run/${RC_SVCNAME}.pid"
command_background=true
output_log="/var/log/benchmark.log"
error_log="/var/log/benchmark.err"
